# systemd/talktask.socket
[Unit]
Description=Socket do projeto de exemplo para 15ยบ Go talks de BH.

[Socket]
ListenStream=80
